# vulnz/src/docker/crontab/mirror
0 0 * * * /mirror.sh >> /var/log/cron_mirror.log 2>&1